# ddns-dnspod.cron.example
#dns record update
*/1 * * * * root /usr/bin/bash /usr/local/ddns-dnspod/ddns-dnspod.sh >>/var/log/ddns-dnspod.log 2>&1
